# academy-hub/codx-room/desktop/supervisord.conf
[program:cairo-dock]
environment=HOME="/home/codx",USER="codx",DISPLAY="%(ENV_DISPLAY)s"
command=cairo-dock
stopsignal=INT
autorestart=true
priority=800
user=codx
redirect_stderr=true
